# dotfiles/config/systemd/user/default.target.wants/syncOrgzly.service
[Unit]
Description=Sync Orgzly through Rclone

[Service]
ExecStart=/home/zarkli/myBin/syncOrgzly
Restart=always
RestartSec=300

[Install]
WantedBy=default.target